# deploy/contentbot.service
[Unit]
Description=Content Factory Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/content-saas-bot
EnvironmentFile=-/opt/content-saas-bot/.env
Environment=PYTHONUNBUFFERED=1
Environment=HTTP_PROXY=
Environment=HTTPS_PROXY=
ExecStart=/opt/content-saas-bot/.venv/bin/python /opt/content-saas-bot/main.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
